# systemd/aws-dynamicdns.timer
[Unit]
Description=Runs the Dynamic DNS service frequently

[Timer]
OnCalendar=*-*-* *:33:00

[Install]
WantedBy=timers.target